# modules/zotero/templates/zotero.service.erb
# NOTE: This file is managed by Puppet
[Unit]
Description="Zotero service"

[Service]
User=zotero
Group=zotero
Environment= LD_LIBRARY_PATH=/usr/lib/xulrunner-24.0/xulrunner/
Restart=always
WorkingDirectory=<%= @base_path %>
ExecStart=/lib64/ld-linux-x86-64.so.2                                   \
 /usr/lib/xulrunner-24.0/xulrunner-sdk/sdk/bin/xpcshell                 \
 -g /usr/lib/xulrunner-24.0/xulrunner                                   \
 -a <%= @base_path %>/translation-server                                \
 -mn                                                                    \
 -v 180 <%= @base_path %>/translation-server/translation-server/init.js \
 >> <%= @log_dir %>/zotero.log 2>&1

[Install]
WantedBy=multi-user.target

